# debian/postinst
#!/bin/bash -x
# postinst script for cuimhne-plugin
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package
VARDIR=/var
USRDIR=/usr
ETCDIR=/etc


case "$1" in
    configure)
        # check that we need to add the user
        if [[ -n $(id -u cuimhne 2>/dev/null) ]]
        then
        	echo "User cuihmne already exists"
        else
            echo "Adding cuihmne user"
            adduser --disabled-login --gecos "" cuimhne
            echo "mus1cM3m0ry\nmus1cM3m0ry" | passwd cuimhne
        fi

        echo "Installing Python libs"
        pip3 install pyudev smbus2 RPi.GPIO rpyc
        
        echo "Preparing directories"
        chown -R cuimhne:cuimhne ${VARDIR}/local/cuimhne
        chmod ug+rwx ${VARDIR}/local/cuimhne
        chmod g+s ${VARDIR}/local/cuimhne
        chmod ug+x ${VARDIR}/local/cuimhne/*.sh
        
        echo "Adding Cuimhne Ceoil scripts to button"
        mkdir -p ${USRDIR}/local/pisound/scripts/pisound-btn/
        ln -sf ${VARDIR}/local/cuimhne/toggle.sh ${USRDIR}/local/pisound/scripts/pisound-btn/toggle_recording.sh 
        ln -sf ${VARDIR}/local/cuimhne/unmount.sh ${USRDIR}/local/pisound/scripts/pisound-btn/unmount_usb.sh 
        ln -sf ${VARDIR}/local/cuimhne/configs/pisound.conf ${ETCDIR}/pisound.conf
        
        echo "Setting up logrotate"
        chmod ug+rwx ${ETCDIR}/chron.hourly/logrotate 
        
        echo "Updating MOTD"
        ln -sf ${VARDIR}/local/cuimhne/motd ${ETCDIR}/motd       
        
        echo "Setting up Nginx"
        ln -sf ${VARDIR}/local/cuimhne/nginx/site/cuimhne ${ETCDIR}/nginx/sites-available/cuimhne
        ln -sf ${ETCDIR}/nginx/sites-available/cuimhne ${ETCDIR}/nginx/sites-enabled/cuimhne
        rm -f ${ETCDIR}/nginx/sites-enabled/default
        htpasswd -bc ${VARDIR}/local/cuimhne/nginx/.htpasswd cuimhne mus1cM3m0ry
        
        echo "Setting up Transmission"
        usermod -a -G cuimhne debian-transmission
        chown -R cuimhne:cuimhne ${VARDIR}/local/opentracker
        chmod ug+rwx ${VARDIR}/local/opentracker
        chmod ug+rx ${VARDIR}/local/opentracker/opentracker
        
        echo "Configuring Systemctl"
        ln -sf ${VARDIR}/local/cuimhne/systemctl/cuimhne.service /usr/lib/systemd/system
        ln -sf ${VARDIR}/local/cuimhne/systemctl/hostname.service /usr/lib/systemd/system
        ln -sf ${VARDIR}/local/cuimhne/systemctl/opentracker.service /usr/lib/systemd/system
        systemctl disable pisound-ctl
        systemctl disable transmission-daemon
        systemctl enable pisound-btn
        systemctl enable hostname
        systemctl enable cuimhne
        systemctl disable ntp
        systemctl disable dnsmasq
        systemctl disable hostapd
        systemctl disable nginx
        systemctl disable jack

        echo "Configuring One Time Scripts"
        touch /boot/ssh
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
